# screenshot/crontab

*/5 * * * * /app/run.sh
@reboot /app/run.sh
